# docker/s6-overlay/s6-rc.d/symlink-frankenphp/up
/etc/s6-overlay/scripts/symlink-frankenphp